# rotation/etc/udev/rules.d/99-rotate-upon-docking.rules
ACTION=="add", SUBSYSTEM=="usb", ATTRS{idVendor}=="06cb", ATTRS{idProduct}=="2819", ENV{DISPLAY}=":0", ENV{XAUTHORITY}="/home/nrndda/.Xauthority", RUN+="/etc/udev/scripts/rotate_normal.sh"
